# meta-cluster-unit/recipes-auto-login/auto-login/auto-login/auto_login_service.service
[Unit]
Description=Auto Login Service
After=systemd-user-sessions.service

[Service]
ExecStart=-/sbin/agetty --autologin root --noclear tty1 $TERM
Type=idle

[Install]
WantedBy=getty.target
